# docs/examples/systemd/aria-underlay-worker.service
[Unit]
Description=Aria Underlay worker daemon
Documentation=docs/runbooks/operator-operations.md
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=aria-underlay
Group=aria-underlay
Environment=RUST_LOG=info
ExecStartPre=/usr/local/bin/aria-underlay-ops check-worker-config --worker-config-path /etc/aria-underlay/worker.json --strict-paths
ExecStart=/usr/local/bin/aria-underlay-worker /etc/aria-underlay/worker.json
StandardOutput=append:/var/log/aria/aria-underlay.log
StandardError=append:/var/log/aria/aria-underlay.log
Restart=on-failure
RestartSec=5s
RuntimeDirectory=aria-underlay
StateDirectory=aria-underlay
LogsDirectory=aria
ReadWritePaths=/var/lib/aria-underlay /var/log/aria /run/aria-underlay
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
CapabilityBoundingSet=
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictRealtime=true
RestrictSUIDSGID=true
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
