# docs/examples/systemd/aria-underlay-product-api.service
[Unit]
Description=Aria Underlay product operations API
Documentation=docs/runbooks/operator-operations.md
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=aria-underlay
Group=aria-underlay
Environment=RUST_LOG=info
Environment=ARIA_UNDERLAY_PRODUCT_API_CONFIG=/etc/aria-underlay/product-api.json
ExecStart=/usr/local/bin/aria-underlay-product-api
StandardOutput=append:/var/log/aria/aria-underlay.log
StandardError=append:/var/log/aria/aria-underlay.log
Restart=on-failure
RestartSec=5s
RuntimeDirectory=aria-underlay
StateDirectory=aria-underlay
LogsDirectory=aria
ReadOnlyPaths=/etc/aria-underlay
ReadWritePaths=/var/lib/aria-underlay /var/log/aria /run/aria-underlay
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
CapabilityBoundingSet=
LockPersonality=true
MemoryDenyWriteExecute=true
RestrictRealtime=true
RestrictSUIDSGID=true
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
